# supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command=/usr/sbin/init

[program:backend]
command=/app/docker_init.sh
